# root/etc/services.d/octoprint/run
#!/usr/bin/execlineb -P
s6-setuidgid abc
octoprint serve --basedir /config
